# Embedded_Linux/07_InitProcess/Busybox_Init/inittab
# Start services when system starts up
::sysinit:/etc/init.d/rcS 3

# start shell
ttyAMA0::askfirst:-/bin/sh

::restart:/sbin/init

# Excute when the system will go down
::shutdown:/etc/init.d/rcS 5

